# automations/crontab
0 4 * * * ~/exevo-pan/automations/maintenance/flushLogs.sh
0 4 * * 3 ~/exevo-pan/automations/maintenance/backupBosses.sh
